# artorize-processor-runner.service
[Unit]
Description=Artorize Image Protection Pipeline Runner
After=network.target

[Service]
Type=simple
User=artorize
Group=artorize
WorkingDirectory=/opt/artorize-processor
Environment="PATH=/opt/artorize-processor/venv/bin"
EnvironmentFile=/opt/artorize-processor/.env

# NOTE: This service is optional. The gateway has internal async workers
# that process jobs from the queue. This runner can be used for:
# - Batch processing workflows
# - Separate processing pool
# - Manual pipeline execution
# To disable: systemctl disable artorize-processor-runner
#
# Pipeline automatically detects and uses GPU if available

ExecStart=/opt/artorize-processor/venv/bin/python -m artorize_runner.pipeline
Restart=on-failure
RestartSec=10
StandardOutput=append:/var/log/artorize/runner.log
StandardError=append:/var/log/artorize/runner-error.log

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/artorize-processor/outputs /opt/artorize-processor/input /opt/artorize-processor/gateway_jobs /var/log/artorize

# Resource limits
LimitNOFILE=65536
TasksMax=4096

[Install]
WantedBy=multi-user.target
